# stix/invstampvest-com-smica83-20260415-b8ec06-9b6c9f1b.json
{
  "type": "bundle",
  "id": "bundle--b8ec06-invstampvest-kimsuky-2026-04-17",
  "objects": [
    {
      "type": "threat-actor",
      "spec_version": "2.1",
      "id": "threat-actor--bea1b40e-0a7a-5778-9e3c-2c46e4480c10",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "name": "Kimsuky",
      "description": "North Korean state-sponsored threat actor also tracked as APT43, Velvet Chollima, Black Banshee, Thallium. Conducts cyber espionage targeting South Korean government, academic, and think-tank organizations. This campaign shows potential geographic expansion with Bulgarian-language invoice lures.",
      "aliases": ["APT43", "Velvet Chollima", "Black Banshee", "Thallium", "ARCHIPELAGO"],
      "threat_actor_types": ["nation-state"],
      "roles": ["agent"],
      "goals": ["espionage", "credential-theft", "intelligence-collection"],
      "sophistication": "expert",
      "resource_level": "government",
      "primary_motivation": "organizational-gain"
    },
    {
      "type": "campaign",
      "spec_version": "2.1",
      "id": "campaign--faktyra-invstampvest-2026-04",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "name": "Kimsuky Invoice Lure Campaign (faktyra/invstampvest)",
      "description": "Campaign using Bulgarian-language invoice-themed ZIP lures (faktyra-850344.zip) delivered via invstampvest.com, with C2 infrastructure at check.nid-log.com. Multi-stage attack chain: CHM dropper -> VBScript stager -> PowerShell keylogger.",
      "first_seen": "2026-02-26T00:00:00.000Z",
      "last_seen": "2026-04-17T00:00:00.000Z"
    },
    {
      "type": "malware",
      "spec_version": "2.1",
      "id": "malware--kimsuky-chm-keylogger-2026",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "name": "Kimsuky CHM Keylogger",
      "description": "Multi-stage malware chain: CHM dropper delivers VBScript recon stager, which establishes persistence via scheduled task and downloads PowerShell keylogger with clipboard monitoring. Uses typo'd User-Agents (Chremo, Edgo) and multipart POST exfiltration.",
      "malware_types": ["keylogger", "spyware", "dropper"],
      "is_family": true,
      "capabilities": ["anti-detection", "persists-after-system-reboot", "captures-input-peripherals", "exfiltrates-data"]
    },
    {
      "type": "indicator",
      "spec_version": "2.1",
      "id": "indicator--invstampvest-domain",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "name": "invstampvest.com",
      "description": "Payload staging domain for Kimsuky invoice lure campaign. Registered 2026-04-14 via Tucows/Winhost.",
      "indicator_types": ["malicious-activity"],
      "pattern": "[domain-name:value = 'invstampvest.com']",
      "pattern_type": "stix",
      "valid_from": "2026-04-14T00:00:00.000Z",
      "kill_chain_phases": [{"kill_chain_name": "mitre-attack", "phase_name": "initial-access"}]
    },
    {
      "type": "indicator",
      "spec_version": "2.1",
      "id": "indicator--nidlog-domain",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "name": "check.nid-log.com",
      "description": "Active Kimsuky C2 server with bootservice.php/checkservice.php/finalservice.php endpoints. Hosted on LightNode VPS at 130.94.29.111.",
      "indicator_types": ["malicious-activity"],
      "pattern": "[domain-name:value = 'nid-log.com' OR domain-name:value = 'check.nid-log.com']",
      "pattern_type": "stix",
      "valid_from": "2026-02-26T00:00:00.000Z",
      "kill_chain_phases": [{"kill_chain_name": "mitre-attack", "phase_name": "command-and-control"}]
    },
    {
      "type": "indicator",
      "spec_version": "2.1",
      "id": "indicator--ip-162-250-72-25",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "name": "162.250.72.25",
      "description": "IP hosting invstampvest.com. WARNING: Shared Winhost hosting - contains legitimate sites. Block domain, not IP.",
      "indicator_types": ["malicious-activity"],
      "pattern": "[ipv4-addr:value = '162.250.72.25']",
      "pattern_type": "stix",
      "valid_from": "2026-04-14T00:00:00.000Z"
    },
    {
      "type": "indicator",
      "spec_version": "2.1",
      "id": "indicator--ip-130-94-29-111",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "name": "130.94.29.111",
      "description": "Active Kimsuky C2 IP on LightNode VPS (AS154177). Hosts check.nid-log.com.",
      "indicator_types": ["malicious-activity"],
      "pattern": "[ipv4-addr:value = '130.94.29.111']",
      "pattern_type": "stix",
      "valid_from": "2026-03-07T00:00:00.000Z",
      "kill_chain_phases": [{"kill_chain_name": "mitre-attack", "phase_name": "command-and-control"}]
    },
    {
      "type": "indicator",
      "spec_version": "2.1",
      "id": "indicator--chm-dropper-hash",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "name": "api_reference.chm SHA256",
      "description": "Kimsuky CHM dropper file. Detection: 10/76 on VirusTotal.",
      "indicator_types": ["malicious-activity"],
      "pattern": "[file:hashes.'SHA-256' = '1eff237dee95172363bfc0342d0389f809f753a6ec5e6848e57b3fd5482e9793']",
      "pattern_type": "stix",
      "valid_from": "2026-04-10T00:00:00.000Z",
      "kill_chain_phases": [{"kill_chain_name": "mitre-attack", "phase_name": "execution"}]
    },
    {
      "type": "indicator",
      "spec_version": "2.1",
      "id": "indicator--vbs-recon-hash",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "name": "payload_1.vbs (recon) SHA256",
      "description": "Kimsuky VBScript recon stager. Detection: 8/76 on VirusTotal.",
      "indicator_types": ["malicious-activity"],
      "pattern": "[file:hashes.'SHA-256' = 'af50f35701916d3909f2727cdcbde1a7af47f46eb8db3996905b1c0725aa133f']",
      "pattern_type": "stix",
      "valid_from": "2026-04-10T00:00:00.000Z",
      "kill_chain_phases": [{"kill_chain_name": "mitre-attack", "phase_name": "execution"}]
    },
    {
      "type": "indicator",
      "spec_version": "2.1",
      "id": "indicator--vbs-stager-hash",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "name": "payload_1.vbs (PS stager) SHA256",
      "description": "Kimsuky VBScript PowerShell bridge stager. Detection: 3/76 on VirusTotal.",
      "indicator_types": ["malicious-activity"],
      "pattern": "[file:hashes.'SHA-256' = 'd7c09e7bf79aa9b786dcd9f870427f4a1110f702646fba9d3835215ad3649d0b']",
      "pattern_type": "stix",
      "valid_from": "2026-04-10T00:00:00.000Z",
      "kill_chain_phases": [{"kill_chain_name": "mitre-attack", "phase_name": "execution"}]
    },
    {
      "type": "indicator",
      "spec_version": "2.1",
      "id": "indicator--persistence-hash",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "name": "Link.ini SHA256",
      "description": "Kimsuky persistence mechanism file. Detection: 7/69 on VirusTotal.",
      "indicator_types": ["malicious-activity"],
      "pattern": "[file:hashes.'SHA-256' = '85f8f8a3f28d2956776fbbd0365cdb78ac8dc1e6ed12818ef18caed0bb2f74c8']",
      "pattern_type": "stix",
      "valid_from": "2026-04-10T00:00:00.000Z",
      "kill_chain_phases": [{"kill_chain_name": "mitre-attack", "phase_name": "persistence"}]
    },
    {
      "type": "indicator",
      "spec_version": "2.1",
      "id": "indicator--ua-chremo",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "name": "Typo'd User-Agent: Chremo",
      "description": "Deliberate Chrome typo in Kimsuky keylogger HTTP exfiltration requests.",
      "indicator_types": ["malicious-activity"],
      "pattern": "[network-traffic:extensions.'http-request-ext'.request_header.'User-Agent' MATCHES 'Chremo/']",
      "pattern_type": "stix",
      "valid_from": "2026-04-10T00:00:00.000Z",
      "kill_chain_phases": [{"kill_chain_name": "mitre-attack", "phase_name": "exfiltration"}]
    },
    {
      "type": "indicator",
      "spec_version": "2.1",
      "id": "indicator--ua-edgo",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "name": "Typo'd User-Agent: Edgo",
      "description": "Deliberate Edge typo in Kimsuky keylogger HTTP exfiltration requests.",
      "indicator_types": ["malicious-activity"],
      "pattern": "[network-traffic:extensions.'http-request-ext'.request_header.'User-Agent' MATCHES 'Edgo/']",
      "pattern_type": "stix",
      "valid_from": "2026-04-10T00:00:00.000Z",
      "kill_chain_phases": [{"kill_chain_name": "mitre-attack", "phase_name": "exfiltration"}]
    },
    {
      "type": "infrastructure",
      "spec_version": "2.1",
      "id": "infrastructure--invstampvest-staging",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "name": "invstampvest.com payload staging",
      "description": "Winhost shared hosting (IIS/10.0, ASP.NET) used for staging invoice-themed malware lures. RapidSSL certificate issued 2026-04-14.",
      "infrastructure_types": ["hosting-malware"]
    },
    {
      "type": "infrastructure",
      "spec_version": "2.1",
      "id": "infrastructure--nidlog-c2",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "name": "nid-log.com C2 infrastructure",
      "description": "Apache/2.4.58 (Win64) PHP/8.2.12 C2 server on LightNode VPS. Health check returns 'Million OK !!!!'. Supports 3-stage payload delivery and keylogger exfiltration.",
      "infrastructure_types": ["command-and-control"]
    },
    {
      "type": "attack-pattern",
      "spec_version": "2.1",
      "id": "attack-pattern--t1566-001",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "name": "Phishing: Spearphishing Attachment",
      "external_references": [{"source_name": "mitre-attack", "external_id": "T1566.001"}]
    },
    {
      "type": "attack-pattern",
      "spec_version": "2.1",
      "id": "attack-pattern--t1056-001",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "name": "Input Capture: Keylogging",
      "external_references": [{"source_name": "mitre-attack", "external_id": "T1056.001"}]
    },
    {
      "type": "attack-pattern",
      "spec_version": "2.1",
      "id": "attack-pattern--t1053",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "name": "Scheduled Task/Job",
      "external_references": [{"source_name": "mitre-attack", "external_id": "T1053"}]
    },
    {
      "type": "relationship",
      "spec_version": "2.1",
      "id": "relationship--kimsuky-uses-campaign",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "relationship_type": "attributed-to",
      "source_ref": "campaign--faktyra-invstampvest-2026-04",
      "target_ref": "threat-actor--bea1b40e-0a7a-5778-9e3c-2c46e4480c10"
    },
    {
      "type": "relationship",
      "spec_version": "2.1",
      "id": "relationship--campaign-uses-malware",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "relationship_type": "uses",
      "source_ref": "campaign--faktyra-invstampvest-2026-04",
      "target_ref": "malware--kimsuky-chm-keylogger-2026"
    },
    {
      "type": "relationship",
      "spec_version": "2.1",
      "id": "relationship--campaign-uses-staging",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "relationship_type": "uses",
      "source_ref": "campaign--faktyra-invstampvest-2026-04",
      "target_ref": "infrastructure--invstampvest-staging"
    },
    {
      "type": "relationship",
      "spec_version": "2.1",
      "id": "relationship--campaign-uses-c2",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "relationship_type": "uses",
      "source_ref": "campaign--faktyra-invstampvest-2026-04",
      "target_ref": "infrastructure--nidlog-c2"
    },
    {
      "type": "relationship",
      "spec_version": "2.1",
      "id": "relationship--malware-indicates-domain",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "relationship_type": "indicates",
      "source_ref": "indicator--invstampvest-domain",
      "target_ref": "campaign--faktyra-invstampvest-2026-04"
    },
    {
      "type": "relationship",
      "spec_version": "2.1",
      "id": "relationship--campaign-uses-phishing",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "relationship_type": "uses",
      "source_ref": "campaign--faktyra-invstampvest-2026-04",
      "target_ref": "attack-pattern--t1566-001"
    },
    {
      "type": "relationship",
      "spec_version": "2.1",
      "id": "relationship--campaign-uses-keylogging",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "relationship_type": "uses",
      "source_ref": "campaign--faktyra-invstampvest-2026-04",
      "target_ref": "attack-pattern--t1056-001"
    },
    {
      "type": "relationship",
      "spec_version": "2.1",
      "id": "relationship--campaign-uses-schtask",
      "created": "2026-04-17T00:00:00.000Z",
      "modified": "2026-04-17T00:00:00.000Z",
      "relationship_type": "uses",
      "source_ref": "campaign--faktyra-invstampvest-2026-04",
      "target_ref": "attack-pattern--t1053"
    }
  ]
}
